# unipi/unipi-default-password/unipicheck.path
[Unit]
Description=Wait for /run/unipi-plc/unipi-id
After=local-fs.target
RequiresMountsFor=/run

[Path]
PathExists=/run/unipi-plc/unipi-id

[Install]
WantedBy=sysinit.target
