# unipi/resize-partition/populate@.service
# Automatically generated by systemd-fstab-generator

[Unit]
Description=Populate created volume %I
DefaultDependencies=no
BindsTo=%i.mount
Conflicts=shutdown.target
After=%i.mount unipi-resize-partition.service
Before=shutdown.target
ConditionPathExists=/%i.cpio.gz

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-/bin/bash -c "gunzip - < /%i.cpio.gz | /usr/bin/cpio -i -D /%I"
ExecStart=/bin/rm -f /%i.cpio.gz
TimeoutSec=0

[Install]
WantedBy=local-fs.target
